# luci-theme-edge/root/etc/uci-defaults/30_luci-theme-edge
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Edge=/luci-static/edge
	commit luci
EOF
exit 0
